# rootfs-overlay/etc/systemd/system/eth-monitor.service
[Unit]
Description=Ethernet Link Monitor (GXL PHY LPA Watchdog)
Documentation=https://github.com/LIyang-linux/B860AV2.1-A
After=network.target fix-gxl-eth.service
ConditionPathExists=/usr/local/bin/eth-monitor.sh

[Service]
Type=simple
ExecStart=/usr/local/bin/eth-monitor.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
